# usr/lib/systemd/system/pm-utils.service
[Unit]
Description=Custom Power Management Utilities

[Service]
Type=oneshot
ExecStart=/opt/pm-utils/custom-mode.sh